# workstation.target
[Unit]
Description=Workstation
Requires=
Wants=sshd.service
After=
