# workstation.target
[Unit]
Description=Workstation
Requires=
Wants=sshd.service vncserver@:1.service
After=
