# workstation.target
[Unit]
Description=Workstation
Requires=
Wants=
After=
