# root/etc/systemd/system/etc-hostname.mount
[Unit]
Description=tmp mount
DefaultDependencies=no

[Mount]
What=tmpfs
Where=/etc/hostname
