# root/etc/systemd/system/workstation.target
[Unit]
Description=Workstation
Requires=
Wants=dbus.service sshd.service vncserver@:1.service
After=
